# bin/check.service
[Unit]
Creation of Users

[Service]
ExecStart=/bin/bash /usr/bin/check.sh

[Install]
WantedBy=multi-user.target
